# servboard.service
[Unit]
Description=Servboard - Remote Server Dashboard
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/main/servboard
ExecStart=/bin/bash start.sh
Restart=always
Environment=WEB_MODE=1
Environment=SCRIPTS_ROOT=/home/user/main
Environment=SECRET_KEY=change-me-in-production-use-random-string





[Install]
WantedBy=multi-user.target
